# ops/com.symphony.report-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.symphony.report-server.plist — Static file server for HTML reports
  Install: cp ops/com.symphony.report-server.plist ~/Library/LaunchAgents/
  Load:    launchctl load ~/Library/LaunchAgents/com.symphony.report-server.plist

  Serves $SYMPHONY_HOME/reports/ on $TOKEN_REPORT_PORT (default 8090).
  KeepAlive: true — restarts automatically if the server crashes.

  SYMPH-131
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.symphony.report-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>http.server</string>
        <string>8090</string>
        <string>--directory</string>
        <string>/Users/clawdilize/.symphony/reports</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/clawdilize/.symphony/reports</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/bin:/usr/sbin:/bin:/sbin</string>
        <key>HOME</key>
        <string>/Users/clawdilize</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/clawdilize/.symphony/logs/report-server-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/clawdilize/.symphony/logs/report-server-stderr.log</string>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
